# etude_de_base/ufwi-administration-suite-ufwi-conf/ufwi_conf/backend/runit/snort_inline/run
#!/bin/sh
exec /usr/bin/snort_inline -m 027 -d -H 1 -Q -l /var/log/snort -c /etc/snort_inline/snort_inline.conf
